# server/supervisor/conf.d/medlib_email_incoming.conf
[program:medlib_queue_email_incoming]
command=php artisan queue:work --daemon --tries=1 --sleep=3 --memory=256 --delay=0 --queue=medlib_email_incoming
directory=/opt/medlib  
stdout_logfile=/var/log/medlib/queue_email_incoming.log  
redirect_stderr=true
user=www-data