# server/crontab
* * * * * root php /var/www/medlib/artisan schedule:run >> /dev/null 2>&1